# .crontab
* * * * * cd /var/www/html/apps.codepolitan.com/public/puller && ./pull.sh